# packages/provisioning/templates/rpi-sentry-health-check.service.j2
[Unit]
Description=Health Check for RPI Sentry

[Service]
ExecStart=/bin/bash /opt/rpi-sentry-health-check/start.sh
Type=oneshot
